# docker/web/supervisor-netspryte.conf
[supervisord]
nodaemon=true

[program:app-nginx]
command=/usr/sbin/nginx

[program:app-uwsgi]
command=/usr/bin/uwsgi --ini /usr/local/share/netspryte/www/uwsgi.ini
