# etc/netspryte.cron
# example cron file for netspryte
MAILTO=root

*/1 * * * * root netspryte-collect-snmp
*/15 * * * * root netspryte-mk-html
